# mgmt/inittab
::respawn:/var/lib/rancher/gen_dnsmasq_conf
::respawn:python /var/lib/rancher/noVNC-0.0.2/utils/websockify --web /var/lib/rancher/noVNC-0.0.2 --cert /var/lib/rancher/self.pem 80 --unix-target /vmmgmt/vm
::restart:/bin/busybox init
